# bsp/board/sdr/rootfs_overlay/etc/systemd/system/mnt-ubootenv.mount
[Unit]
Description=Mount partition with uboot env

[Mount]
What=/dev/mmcblk0p1
Where=/mnt/ubootenv/
Type=vfat
Options=defaults

[Install]
WantedBy=multi-user.target

